# etc/services.d/xvfb/run
#!/usr/bin/with-contenv /bin/bash
RES=${RESOLUTION:-1280x800x24}
s6-setuidgid kali Xvfb ":0" -screen 0 "$RES"
