# etc/services.d/kde/run
#!/usr/bin/with-contenv /bin/bash
export DISPLAY=":0"
sudo -u kali startplasma-x11
